# deploy/reindex.timer
[Unit]
Description=Daily rebuild of the event photo face-embeddings index

[Timer]
OnCalendar=*-*-* 04:00:00
Persistent=true

[Install]
WantedBy=timers.target
